# configs/supervisor/call_me_ai_telegram_polling.conf
[program:call_me_ai_telegram_polling]
directory=/home/idurnev/apps/call_me_ai
user=idurnev
group=idurnev
command=/bin/bash -lc 'set -a; source /home/idurnev/apps/call_me_ai/.env; set +a; exec /home/idurnev/apps/call_me_ai/venv/bin/python /home/idurnev/apps/call_me_ai/bot_polling.py'
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/call_me_ai_telegram_polling.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/call_me_ai_telegram_polling.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1"
